# files/system/usr/lib/systemd/system/xremap.service
[Unit]
Description=Xremap
After=sddm-boot.service

[Service]
ExecStart=/usr/bin/xremap /usr/etc/xremap.yml

[Install]
WantedBy=multi-user.target
