# files/system/etc/systemd/system/flatpak_run_script.timer
[Unit]
Description=Automatically creates run scrips for flatpak apps.

[Timer]
OnBootSec=10s
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
